# container/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile=/denodo/supervisord.log
pidfile=/var/run/supervisord.pid                ; pidfile location
nodaemon=true                                   ; run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=root                                       ; default user
childlogdir=/denodo                             ; where child log files will live


[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:denodoinit]
priority=100
command=/opt/denodo/denodoinstall.sh
process_name=denodoinit
autostart=true
autorestart=false
startinorder=true
startnextafter=exited
user=root

[program:vqlserver]
priority=200
command=/denodo/bin/vqlserver_startup.sh
process_name=vqlserver
autorestart=false
startretries=0
user=root
startinorder=true
startnextafter=exited

[program:datacatalog]
priority=300
command=/denodo/bin/datacatalog_startup.sh
process_name=datacatalog
autorestart=false
startretries=0
user=root
startinorder=true
startnextafter=exited

[program:datamonitoring]
priority=400
command=/denodo/bin/diagnosticmonitoringtool_startup.sh
process_name=datacatalog
autorestart=false
startretries=0
user=root
startinorder=true
startnextafter=exited

[./program:monitoring]
priority=500
command=/denodo/tools/monitor/denodo-monitor/bin/denodomonitor_startup.sh
process_name=datacatalog
autorestart=false
startretries=0
user=root
startinorder=true
startnextafter=exited
